# systemd/innotech_scraper_api.service
[Unit]                                            
Description=Innotech Scraper FastAPI Server       
After=network.target

[Service]      
User=innotech_scraper
WorkingDirectory=/opt/innotech_scraper
ExecStart=/opt/innotech_scraper/selenium-env/bin/uvicorn innotechscraper_api:app --host 0.0.0.0 --port 8001 --no-access-log --log-level error
Restart=always

[Install]
WantedBy=multi-user.target
